# uevent/diskMount/m1/99-udisk.rules
KERNEL=="sd[a-z][0-9]", SUBSYSTEM=="block", RUN+="/home/pangwz/test/uevent/diskMount/udev_disk_auto_mount.sh %k %N"
