# udev/info/99-info-rowan.rules
RUN+="/home/pangwz/test/uevent/info/udev_info.sh %k %n"
